# systemd/x-collector.timer
[Unit]
Description=Run X Web3 Digest collector hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
